# config/com.road4ai.mimo-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.road4ai.mimo-server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/llama-server</string>
        <string>-m</string>
        <string>/Users/shagwu/.local/share/models/mimo-7b-rl/MiMo-7B-RL-Q4_K_M.gguf</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8000</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/shagwu/.local/share/models/mimo-server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shagwu/.local/share/models/mimo-server-error.log</string>
</dict>
</plist>
